# files/systemd/user/distrobox-upgrade.timer
[Unit]
Description=Run Distrobox upgrade daily

[Timer]
OnCalendar=daily
RandomizedDelaySec=2h
Persistent=true

[Install]
WantedBy=timers.target
