# files/systemd/user/distrobox-upgrade.service
[Unit]
Description=Upgrade all Distrobox containers

[Service]
Type=oneshot
ExecStart=%h/.local/bin/distrobox-upgrade-all
